# imaporter/rootfs/etc/services.d/spamd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start SpamAssassin daemon (spamd)
# Only starts if spamassassin_enabled is true in app config
# ==============================================================================

if bashio::config.true 'spamassassin_enabled'; then
    # Build local.cf from HA options
    SCORE=$(bashio::config 'spamassassin_required_score')
    bashio::log.info "Configuring SpamAssassin (required_score: ${SCORE})..."
    echo "required_score ${SCORE}" > /etc/mail/spamassassin/local.cf

    # Append user-supplied rules from /config/imaporter/spamassassin.cf if present
    EXTERNAL_CF="/homeassistant/imaporter/spamassassin.cf"
    if [ -f "${EXTERNAL_CF}" ]; then
        bashio::log.info "Appending custom SpamAssassin rules from ${EXTERNAL_CF}..."
        cat "${EXTERNAL_CF}" >> /etc/mail/spamassassin/local.cf
    fi

    bashio::log.info "Updating SpamAssassin rules..."
    if sa-update; then
        bashio::log.info "SpamAssassin rules updated successfully."
    else
        bashio::log.warning "SpamAssassin rule update failed or no new updates available."
    fi
    bashio::log.info "Starting SpamAssassin daemon..."
    exec spamd --listen=127.0.0.1 --allow-tell --nouser-config
else
    bashio::log.info "SpamAssassin disabled in configuration, skipping."
    exec sleep infinity
fi
