# imaporter.service
[Unit]
Description=IMAPorter (Account A -> Gmail)
After=network-online.target spamassassin.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/imaporter/venv/bin/python /opt/imaporter/imaporter.py
WorkingDirectory=/opt/imaporter
Restart=always
RestartSec=30
User=imaporter
LoadCredential=source_password:/etc/imaporter/source.secret
LoadCredential=destination_password:/etc/imaporter/destination.secret
ProtectHome=yes

[Install]
WantedBy=multi-user.target
