# imaporter/rootfs/etc/services.d/imaporter/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start IMAPorter service
# Generates config.ini from HA add-on options and launches imaporter.py
# ==============================================================================

CONFIG="/data/config.ini"

bashio::log.info "Generating IMAPorter configuration from add-on options..."

# Build config.ini from HA options.json
cat > "${CONFIG}" <<EOF
[source]
host = $(bashio::config 'source_host')
port = $(bashio::config 'source_port')
username = $(bashio::config 'source_username')
password = $(bashio::config 'source_password')
ssl = $(bashio::config 'source_ssl')
folder = $(bashio::config 'source_folder')
delete_on_source = $(bashio::config 'source_delete')

[destination]
host = $(bashio::config 'destination_host')
port = $(bashio::config 'destination_port')
username = $(bashio::config 'destination_username')
password = $(bashio::config 'destination_password')
ssl = $(bashio::config 'destination_ssl')
ham_folder = $(bashio::config 'destination_ham_folder')
ham_label = $(bashio::config 'destination_ham_label')

[spamassassin]
enabled = $(bashio::config 'spamassassin_enabled')
max_size = $(bashio::config 'spamassassin_max_size')
EOF

LOG_LEVEL=$(bashio::config 'log_level')

bashio::log.info "Starting IMAPorter (log level: ${LOG_LEVEL})..."
exec python3 /opt/imaporter/imaporter.py \
    --config "${CONFIG}" \
    --log-level "${LOG_LEVEL}"
