# imaporter/rootfs/etc/services.d/imaporter/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start IMAPorter service
# Generates config.ini from HA add-on options and launches imaporter.py
# ==============================================================================

CONFIG="/data/config.ini"

bashio::log.info "Generating IMAPorter configuration from add-on options..."

# Build config.ini from HA options.json
> "${CONFIG}"

for source in $(bashio::config 'sources|keys'); do
    echo "[source_${source}]" >> "${CONFIG}"
    echo "host = $(bashio::config "sources[${source}].host")" >> "${CONFIG}"
    echo "port = $(bashio::config "sources[${source}].port")" >> "${CONFIG}"
    echo "username = $(bashio::config "sources[${source}].username")" >> "${CONFIG}"
    echo "password = $(bashio::config "sources[${source}].password")" >> "${CONFIG}"
    echo "ssl = $(bashio::config "sources[${source}].ssl")" >> "${CONFIG}"
    echo "folder = $(bashio::config "sources[${source}].folder")" >> "${CONFIG}"
    echo "delete_on_source = $(bashio::config "sources[${source}].delete")" >> "${CONFIG}"
    echo "ham_label = $(bashio::config "sources[${source}].ham_label")" >> "${CONFIG}"
    echo "" >> "${CONFIG}"
done

cat >> "${CONFIG}" <<EOF
[destination]
host = $(bashio::config 'destination_host')
port = $(bashio::config 'destination_port')
username = $(bashio::config 'destination_username')
password = $(bashio::config 'destination_password')
ssl = $(bashio::config 'destination_ssl')
ham_folder = $(bashio::config 'destination_ham_folder')

[spamassassin]
enabled = $(bashio::config 'spamassassin_enabled')
max_size = $(bashio::config 'spamassassin_max_size')
EOF

LOG_LEVEL=$(bashio::config 'log_level')

bashio::log.info "Starting IMAPorter (log level: ${LOG_LEVEL})..."
exec python3 /opt/imaporter/imaporter.py \
    --config "${CONFIG}" \
    --log-level "${LOG_LEVEL}"
